# buildpack/profile.d/10-docmost.sh
# Docmost computed envs
export DOCMOST_VERSION
